# simplerisk-compose/web/common/etc/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
childlogdir=/var/log/supervisor
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisor/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:httpd]
command=/etc/apache2/foreground.sh
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/httpd.log
stderr_logfile=/var/log/supervisor/httpd.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stopsignal=TERM

[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/cron.log
stderr_logfile=/var/log/supervisor/cron.err
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stopsignal=TERM
