# simplerisk-minimal/common/etc/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                   ; (the path to the socket file)
chmod=0700

[supervisord]
childlogdir=/tmp
logfile=/dev/stdout                         ; (log to stdout for container-native logging)
logfile_maxbytes=0                          ; (disable rotation — stdout is not seekable)
loglevel=info                               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid                ; (use /tmp — always writable in containers)
nodaemon=true                               ; (start in foreground if true;default false)
minfds=1024                                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=10

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
